# systemd/seasoning-ticker.service
# ~/.config/systemd/user/seasoning-ticker.service
[Unit]
Description=Seasoning second-aligned ticker

[Service]
ExecStart=/usr/bin/seasoning ticker 6
Restart=always

[Install]
WantedBy=default.target
